# deploy/systemd/aspida-hura.service.d/prefix-cache.conf
[Service]
Environment=ASPIDA_PREFIX_CACHE=1
Environment=GGML_CUDA_DISABLE_GRAPHS=1
